# laptop/rc.local
#!/bin/bash
echo D015 > /proc/acpi/wakeup

exit 0
